# .config/systemd/user/set-earth-background.service
[Unit]
Description=Live Earth Wallpaper

[Service]
Type=simple
ExecStart=/home/remo/.local/bin/set_earth_background
Nice=19
